# roles/init-idc-centos7/files/ecs_mq/src/62-ecs-mq.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="eth*", RUN+="/sbin/ecs_mq_rps_rfs"
